# nix/hm/litellm/deploy/litellm.service
# LiteLLM System Service
# Install: sudo cp litellm.service /etc/systemd/system/
# Enable: sudo systemctl enable litellm
# Start: sudo systemctl start litellm
# Logs: sudo journalctl -u litellm -f
#
# IMPORTANT: This is a system service (not user service) to survive SSH disconnections

[Unit]
Description=LiteLLM Proxy Server
Documentation=https://docs.litellm.ai/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10

# Run as litellm user (created during setup)
User=litellm
Group=litellm

# Working directory
WorkingDirectory=/opt/litellm

# Environment file for secrets (optional, secrets embedded in config.yaml)
EnvironmentFile=-/etc/litellm/secrets.env

# Environment variables
Environment="LITELLM_LOG=INFO"
Environment="PYTHONUNBUFFERED=1"
Environment="UV_DEFAULT_INDEX=https://mirrors.tuna.tsinghua.edu.cn/pypi/web/simple"
Environment="UV_PYTHON_INSTALL_MIRROR=https://gh-proxy.com/github.com/astral-sh/python-build-standalone/releases/download"

# Run LiteLLM using uvx (automatically manages Python environment)
# uvx is copied to /usr/local/bin by deploy.sh
ExecStart=/usr/local/bin/uvx --python 3.11 --with 'litellm[proxy]==1.80.11' --with 'httpx[socks]' \
    litellm==1.80.11 --config /etc/litellm/config.yaml --host 0.0.0.0 --port 4000

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/litellm

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
